# systemd/cdc-emitter.service
[Unit]
Description=CDC Binlog Emitter Service
Documentation=https://github.com/FrontX-Vicky/binlogTriggers
After=network.target mysql.service redis-server.service
Wants=mysql.service redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/go-workspace/mysql_changelog_publisher
ExecStart=/var/www/go-workspace/mysql_changelog_publisher/cdc-emitter
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cdc-emitter

# Environment
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
