# systemd/cdc-subscribers.service
[Unit]
Description=CDC Subscribers Service
Documentation=https://github.com/FrontX-Vicky/binlogTriggers
After=network.target redis.service cdc-emitter.service
Wants=redis.service

[Service]
Type=simple
User=cdc
Group=cdc
WorkingDirectory=/opt/cdc-platform
Environment="ENV_FILES=/opt/cdc-platform/.env.lead_events,/opt/cdc-platform/.env.user_events"
ExecStart=/opt/cdc-platform/bin/subscribers
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cdc-subscribers

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cdc-platform/logs

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Health monitoring
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
